# software/etc/systemd/system/mainController.service
[Unit]
Description=Vehicle Immobilization Device

[Service]
Type=forking
ExecStart=/usr/bin/mainController.sh

[Install]
WantedBy=multi-user.target
